# tatu-bots/deploy/tatulms-watchdog.timer
[Unit]
Description=tatulmsbot watchdog — har 2 daqiqada tekshiradi

[Timer]
OnBootSec=3min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
